# medical-bot.service
[Unit]
Description=Medical Reports Telegram Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=botuser
WorkingDirectory=/home/botuser/medical-bot
Environment=PATH=/home/botuser/medical-bot/venv/bin
Environment=PYTHONPATH=/home/botuser/medical-bot
ExecStart=/home/botuser/medical-bot/venv/bin/python app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=medical-bot

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/botuser/medical-bot /home/botuser/backups
ProtectKernelTunables=yes
ProtectControlGroups=yes

# Resource limits
MemoryLimit=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target

