# medical-bot.service
[Unit]
Description=Medical Telegram Bot
After=network.target

[Service]
User=botuser
WorkingDirectory=/home/botuser/medical-bot
EnvironmentFile=-/home/botuser/medical-bot/config.env
ExecStart=/home/botuser/medical-bot/venv/bin/python app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
